# packages/timelapse/timelapse-press-button@.service
[Unit]
DefaultDependencies=no
Before=shutdown.target
Requires=pisugar-server.service


[Service]
Type=oneshot
User=timelapse
ExecStart=timelapse press-button %i

[Install]
WantedBy=multi-user.target