# packages/timelapse/timelapse-hotspot.service
[Service]
Type=simple
ExecStart=timelapse hotspot