# packages/timelapse/timelapse-website.service
[Service]
Type=oneshot
User=timelapse
ExecStart=timelapse website