# packages/timelapse/timelapse-website.service
[Service]
User=timelapse
ExecStart=timelapse-website --port=80
AmbientCapabilities=CAP_NET_BIND_SERVICE