# packages/timelapse/timelapse-website.service
[Service]
User=timelapse
ExecStart=timelapse website